# AllskyMQTTBatteryMonitor.service
[Unit]
Description=AllSky MQTT Battery Monitor
After=network.target
Wants=network.target

[Service]
Type=simple
User=asi
Group=asi
WorkingDirectory=/home/asi/AllskyMQTTBatteryMonitor
Environment=PATH=/home/asi/AllskyMQTTBatteryMonitor/.venv/bin
ExecStart=/home/asi/AllskyMQTTBatteryMonitor/.venv/bin/python /home/asi/AllskyMQTTBatteryMonitor/main.py --broker 172.17.204.35 --topic solar/pmallsky/data --filepathprogress /home/asi/AllskyMQTTBatteryMonitor/BatteryPercent.txt --filepathstatic /home/asi/AllskyMQTTBatteryMonitor/StaticVars.txt
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
